# files/etc/crontabs/root
0 */2 * * * /sbin/ifdown wan6; sleep 5; /sbin/ifup wan6
